# deploy/doa-md-poll.service
# /etc/systemd/system/doa-md-poll.service
# 1회성 폴링. 별도 timer 가 호출.

[Unit]
Description=Doa MD Agent — poll once (auto_login + crawl + dump + notify)
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/doa-md-agent
EnvironmentFile=/home/ubuntu/doa-md-agent/.env
ExecStart=/home/ubuntu/.local/bin/uv run python -m crawler.auto_login
ExecStart=/home/ubuntu/.local/bin/uv run python -m crawler.run crawl
ExecStart=/home/ubuntu/.local/bin/uv run python -m crawler.run sales-all
ExecStart=/home/ubuntu/.local/bin/uv run python -m crawler.run dump-json
ExecStart=/home/ubuntu/.local/bin/uv run python -m crawler.notify
StandardOutput=append:/home/ubuntu/doa-md-agent/data/poll.log
StandardError=append:/home/ubuntu/doa-md-agent/data/poll.log
